# deployment/forestguard.service
[Unit]
Description=ForestGuard API - Wildfire Recovery System
After=network.target

[Service]
Type=simple
User=opc

WorkingDirectory=/home/opc/wildfire-recovery-argentina
EnvironmentFile=/home/opc/wildfire-recovery-argentina/.env
ExecStart=/home/opc/wildfire-recovery-argentina/venv/bin/gunicorn app.main:app \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --timeout 120

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging (handled by journald)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forestguard

[Install]
WantedBy=multi-user.target
